# deploy/freebsd/trmnl_hsl
#!/bin/sh
#
# FreeBSD rc.d service for the TRMNL HSL departure board server.
#
# PROVIDE: trmnl_hsl
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# rc.subr reads the contract variables set below by name, so they look
# unused to shellcheck.
# shellcheck disable=SC2034
#
# shellcheck source=/dev/null
. /etc/rc.subr

name="trmnl_hsl"
rcvar="trmnl_hsl_enable"

: "${trmnl_hsl_user:=www}"
: "${trmnl_hsl_appdir:=/opt/trmnl-hsl}"
: "${trmnl_hsl_env:=}"

pidfile="/var/run/${name}.pid"
logfile="/var/log/${name}.log"

command="/usr/sbin/daemon"
command_args="-f -r -P ${pidfile} -t ${name} -o ${logfile} /usr/bin/env ${trmnl_hsl_env} ${trmnl_hsl_appdir}/run.sh"
procname="/usr/sbin/daemon"

start_precmd="trmnl_hsl_precmd"

trmnl_hsl_precmd() {
  # Pre-create PID/log owned by the run user (their dirs are root-only).
  install -o "${trmnl_hsl_user}" -g "${trmnl_hsl_user}" -m 644 /dev/null "${pidfile}"
  install -o "${trmnl_hsl_user}" -g "${trmnl_hsl_user}" -m 644 /dev/null "${logfile}"
}

load_rc_config $name
run_rc_command "$1"
